# testdata/crontab/invalid.cron
# Crontab with invalid entries for testing error handling

# Valid entry
0 0 * * * /usr/bin/valid.sh

# Invalid: too few fields
0 0 * *

# Invalid: minute out of range
60 0 * * * /usr/bin/invalid-minute.sh

# Invalid: hour out of range
0 25 * * * /usr/bin/invalid-hour.sh

# Valid entry
@daily /usr/bin/another-valid.sh

# Invalid: garbage input
this is not a cron entry at all

# Valid entry with comment
*/5 * * * * /usr/bin/frequent.sh # Every 5 minutes

